# fw/linux/petalinux_emmc_spec/meta-user/recipes-apps/mcoiapp/files/mcoiapp.service
[Unit]
Description=mcoiapp

[Service]
ExecStart=/usr/bin/mcoiapp
StandardOutput=journal+console

[Install]
WantedBy=multi-user.target

